# .config/systemd/user/battery-notify.service
[Unit]
Description=Battery level checker with desktop notifications

[Service]
Type=oneshot
ExecStart=%h/.config/battery-notify/check.sh
