# .config/systemd/user/gdrive-backup.service
[Unit]
Description=Monthly Google Drive Backup

[Service]
Type=oneshot
ExecStart=/home/ziad/Documents/backup/gdrive_backup.sh
