# systemd/user/battery-notify.timer
[Unit]
Description=Run battery checker every minute

[Timer]
OnBootSec=30
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
